# .config/autostart/camera_app.desktop
[Desktop Entry]
Type=Application
Name=Camera App
Exec=/home/pi/venv/bin/python /home/pi/camera_app.py
StartupNotify=false
