# .config/autostart/camera_app.desktop
[Desktop Entry]
Type=Application
Name=Camera App
Exec=python3 /home/kev/pikon_ii/camera_app.py
StartupNotify=false
